# src/tpl/systemd.service
[Unit]
Description=CINV backend

[Service]
ExecStart=CINV_START

[Install]
WantedBy=multi-user.target
